# compose/stacks/h-yt-dlp/scripts/crontab
# 04:00 todos os dias — sincroniza playlists do YouTube com a biblioteca.
# Output capturado por supercronic em stdout (visível via `docker logs`).
0 4 * * * /usr/local/bin/sync.sh
